# media-sound/RoonBridge/files/roonbridge.init.d.rt
#!/sbin/openrc-run

export SSD_NICELEVEL="-20"

command="chrt"
command_args="
	-r 51 /opt/RoonBridge/Bridge/RoonBridge
"

command_background=yes
pidfile=/run/roon.pid

depend() {
    use alsasound
    after bootmisc
}
