# media-sound/RoonBridge/files/roonbridge.init.d.rt
#!/sbin/openrc-run

user="root:root"
logfile="/tmp/roon.log"

command="chrt"
command_args="
        -r 51 /opt/RoonBridge/start.sh -f $logfile
"

command_background=yes
pidfile=/run/roon.pid

#need net
depend() {
    use alsasound
    after bootmisc
}

start_pre() {
    checkpath --file --owner $user --mode 0644 $logfile
}
